# supervisor/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=warn

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stdout_logfile=/var/log/nginx-stdout.log
stderr_logfile=/var/log/nginx/error.log
priority=10

[program:alist]
command=/app/alist server --data /data/alist
autostart=true
autorestart=true
stdout_logfile=/var/log/alist.log
stderr_logfile=/var/log/alist.log
priority=20

[program:rclone]
command=/usr/bin/rclone rcd --rc-addr=127.0.0.1:5572 --rc-no-auth --config=/data/rclone/rclone.conf --cache-dir=/data/rclone/cache
autostart=true
autorestart=true
stdout_logfile=/var/log/rclone.log
stderr_logfile=/var/log/rclone.log
priority=30

[program:warp]
command=/usr/bin/warp-plus --bind 127.0.0.1:1080
autostart=true
autorestart=true
stdout_logfile=/var/log/warp.log
stderr_logfile=/var/log/warp.log
priority=35

[program:rclone-proxy]
command=/usr/bin/rclone rcd --rc-addr=127.0.0.1:5573 --rc-no-auth --config=/data/rclone/rclone.conf --cache-dir=/data/rclone/cache-proxy
autostart=true
autorestart=true
environment=ALL_PROXY="socks5://127.0.0.1:1080"
stdout_logfile=/var/log/rclone-proxy.log
stderr_logfile=/var/log/rclone-proxy.log
priority=38

[program:api]
command=node /app/server/index.js
autostart=true
autorestart=true
stdout_logfile=/var/log/api.log
stderr_logfile=/var/log/api.log
environment=NODE_ENV="production"
priority=40

[program:autosync]
command=/bin/bash -c 'if [ -z "$SYNC_DEST" ]; then sleep infinity; fi; echo "[AutoSync] Started syncing /data to $SYNC_DEST every 5 minutes"; while true; do sleep 300; /usr/bin/rclone sync /data "$SYNC_DEST" --exclude "rclone/cache/**" --exclude "alist/data/temp/**" --exclude "alist/temp/**" --exclude "alist/data/bleve/**" -v; done'
autostart=true
autorestart=true
stdout_logfile=/var/log/autosync.log
stderr_logfile=/var/log/autosync.log
priority=50
